# com.nudge.launcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Nudge launchd agent.

  This file is bundled inside the app at Contents/Resources/com.nudge.launcher.plist.
  On first launch, Nudge copies it to ~/Library/LaunchAgents/com.nudge.launcher.plist
  with {{BINARY_PATH}} replaced by the actual app binary path, then loads it via:
      launchctl load ~/Library/LaunchAgents/com.nudge.launcher.plist

  RunAtLoad=true launches Nudge at login. KeepAlive is intentionally off so the
  user can Cmd+Q to quit during the session without launchd immediately respawning.

  TO UNINSTALL MANUALLY:
      launchctl unload ~/Library/LaunchAgents/com.nudge.launcher.plist
      rm ~/Library/LaunchAgents/com.nudge.launcher.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nudge.launcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{BINARY_PATH}}</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
